# scripts/setup/configs/friendlyelec/rk3588-udev.rules.tmpl
# FriendlyElec RK3588/RK3588S udev rules
# Hardware-specific device configuration and optimization

# I/O scheduler optimizations for RK3588/RK3588S storage devices
# NVMe drives - use mq-deadline for better performance
ACTION=="add|change", KERNEL=="nvme[0-9]*", ATTR{queue/scheduler}="mq-deadline"

# eMMC - use deadline scheduler for better sequential performance
ACTION=="add|change", KERNEL=="mmcblk[0-9]*", ATTR{queue/scheduler}="deadline"

# SD cards - use noop for minimal overhead
ACTION=="add|change", KERNEL=="mmcblk[0-9]*p*", ATTR{queue/scheduler}="noop"

# Set read-ahead for storage devices
ACTION=="add|change", KERNEL=="nvme[0-9]*", ATTR{bdi/read_ahead_kb}="512"
ACTION=="add|change", KERNEL=="mmcblk[0-9]*", ATTR{bdi/read_ahead_kb}="256"

# VPU (Video Processing Unit) device permissions
KERNEL=="mpp_service", GROUP="video", MODE="0660", TAG+="uaccess"

# DRM devices for hardware acceleration
KERNEL=="card[0-9]*", GROUP="video", MODE="0664", TAG+="uaccess"
KERNEL=="renderD[0-9]*", GROUP="render", MODE="0664", TAG+="uaccess"

# GPIO devices for hardware access
SUBSYSTEM=="gpio", GROUP="gpio", MODE="0664"
KERNEL=="gpiochip*", GROUP="gpio", MODE="0664"

# I2C devices for sensor access
KERNEL=="i2c-[0-9]*", GROUP="i2c", MODE="0664"

# SPI devices
KERNEL=="spidev*", GROUP="spi", MODE="0664"

# PWM devices for fan control
KERNEL=="pwm*", GROUP="gpio", MODE="0664"

# RTC device permissions
KERNEL=="rtc[0-9]*", GROUP="audio", MODE="0664"

# USB devices optimization
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="2207", ATTR{idProduct}=="*", MODE="0664", GROUP="plugdev"

# Network interface naming for FriendlyElec devices
SUBSYSTEM=="net", ACTION=="add", DRIVERS=="?*", ATTR{address}=="*", ATTR{dev_id}=="0x0", ATTR{type}=="1", NAME="eth0"

# Set CPU governor on boot (backup method)
ACTION=="add", SUBSYSTEM=="cpu", KERNEL=="cpu[0-9]*", RUN+="/bin/sh -c 'echo performance > /sys/devices/system/cpu/cpu%n/cpufreq/scaling_governor'"
